# better_dontforget/systemd_units/better-dontforget-notify.service
[Unit]
Description=Better Dontforget pending reminder notifications
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=__BINARY__ notify-pending

[Install]
WantedBy=default.target
